# deploy/ansible/roles/wpcomposer/templates/wpcomposer-pipeline.timer.j2
[Unit]
Description=Run WP Composer pipeline hourly

[Timer]
OnCalendar=hourly
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
